# deployment/systemd/lsx-backup-db.service
[Unit]
Description=LernsystemX PostgreSQL Database Backup
Documentation=https://github.com/your-org/lernsystemx
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres

# Load environment variables from production config
EnvironmentFile=/opt/lsx/backend/.env.production

# Execute backup script
ExecStart=/bin/bash /opt/lsx/deployment/backup/backup_db.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lsx-backup-db

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/lsx

# Timeout (max 30 minutes for backup)
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
